# crontab
# run every day at 02:01am
01 10 * * * cd /opt/apps && ./exec_crypto.sh >> /var/log/crypto.log 2>&1
#* * * * * cd /opt/apps && ./exec_crypto.sh >> /var/log/crypto.log 2>&1
